# chrome-debug.service
[Unit]
Description=Chrome with Remote Debugging for SMS Scraping
After=network.target
Wants=network.target

[Service]
Type=forking
User=root
WorkingDirectory=/var/www/html/telegram/teleram-otp-bot
ExecStartPre=/bin/bash -c 'pkill -f "chrome.*remote-debugging" || true'
ExecStartPre=/bin/sleep 2
ExecStart=/bin/bash -c '/var/www/html/telegram/teleram-otp-bot/start_chrome_debug_linux.sh'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chrome-debug

# Cleanup on stop
ExecStop=/bin/bash -c 'pkill -f "chrome.*remote-debugging" || true'
ExecStopPost=/bin/sleep 2

[Install]
WantedBy=multi-user.target

